# templates/redis.service.erb
[Unit]
Description=Redis In-Memory Data Store
After=network.target

[Service]
User=<%= @redis_user %>
Group=<%= @redis_group%>
ExecStart=/usr/local/bin/redis-server <%= @conf %>
ExecStop=/usr/local/bin/redis-cli -p <%= @port %><%= " -a #{@requirepass}" unless @requirepass.empty? %> shutdown
Restart=always

[Install]
WantedBy=multi-user.target
